# cfg-templates/crontab
30 * * * * /scripts/topology_submit
45 * * * * /scripts/topology_update
